# templates/etc/system.d/system/docker.service.d/http-proxy.conf.j2
[Service]
Environment="HTTP_PROXY={{ docker_proxy }}"

